# Surport_Files/supervisord.conf
[program:shadowsocks]
command=sslocal -c ~/.script/Surport_Files/shadowsocks.json
autostart=true 
autorestart=true
user=root
log_stderr=true
logfile=/var/log/shadowsocks.log
